# files/podman-junker
#!/sbin/openrc-run

podman_user="junker"

source "/etc/init.d/podman"
